# deploy/supervisord/dog-ai.conf
[program:dog-ai]
command=/opt/dog-ai/venv/bin/python /opt/dog-ai/pipeline.py
directory=/opt/dog-ai
user=dog-ai
autostart=true
autorestart=true
startsecs=5
startretries=10
stopwaitsecs=30
stopsignal=INT
stdout_logfile=/var/log/dog-ai/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/dog-ai/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=TOGETHER_API_KEY="%(ENV_TOGETHER_API_KEY)s",HF_TOKEN="%(ENV_HF_TOKEN)s"
