# deploy/supervisord/karaos.conf
[program:karaos]
command=/opt/karaos/venv/bin/python /opt/karaos/pipeline.py
directory=/opt/karaos
user=karaos
autostart=true
autorestart=true
startsecs=5
startretries=10
stopwaitsecs=30
stopsignal=INT
stdout_logfile=/var/log/karaos/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/karaos/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
; Only the REQUIRED var is interpolated here: supervisord expands %(ENV_X)s at
; config-read time and ERRORS if the var is unset — naming optional vars here
; would make them structurally required. Optional vars (HF_TOKEN, GROQ_API_KEY,
; TAVILY_API_KEY) reach the pipeline by plain inheritance from supervisord's
; own environment (export them before starting supervisord).
environment=TOGETHER_API_KEY="%(ENV_TOGETHER_API_KEY)s"
