# deploy/systemd/karaos.service
[Unit]
Description=karaos cognitive runtime
Documentation=https://github.com/HungryFingerss/KaraOS
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=karaos
WorkingDirectory=/opt/karaos
ExecStart=/opt/karaos/venv/bin/python /opt/karaos/pipeline.py
EnvironmentFile=/etc/karaos/karaos.env
Restart=on-failure
RestartSec=5s
StartLimitBurst=5
StartLimitIntervalSec=60s
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
